# ops/systemd/marketing-compose.service
[Unit]
Description=Carolina Growth docker compose stack
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/marketing
ExecStart=/usr/bin/docker compose -f /opt/marketing/docker-compose.prod.yml --env-file /opt/marketing/.env.prod up -d
ExecStop=/usr/bin/docker compose -f /opt/marketing/docker-compose.prod.yml --env-file /opt/marketing/.env.prod down
TimeoutStartSec=0

[Install]
WantedBy=multi-user.target
